# systemd/blocktimer.service.in
[Unit]
Description=Blocktimer ip block scheduler
After=network.target

[Service]
ExecStart=@BLOCKTIMER_BIN@
# Restart=on-failure

[Install]
WantedBy=multi-user.target
